# environments/laravel-queue-worker/crontab
* * * * * php /var/www/portal/artisan schedule:run >> /dev/null 2>&1
